# autoshutdown.timer
[Unit]
Description=periodic check for running program

[Timer]
OnBootSec=300
OnUnitActiveSec=10
